# dcron.supervisor.conf
[program:dcron]
command=/bin/bash -c "cd /app/dcrontab && nice -n 5 ./dcron"
redirect_stderr=true
