# dcron.supervisor.conf
[program:dcron]
command=/bin/bash -c "cd /app/dcrontab && [ -z '$NODEID' ] && nice -n 5 ./dcron -nodeid 1 || nice -n 5 ./dcron -nodeid $NODEID"
redirect_stderr=true
